# deploy/xpo-server.service
[Unit]
Description=xpo.sh tunnel server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=xpo
Group=xpo
ExecStart=/usr/local/bin/xpo-server
Restart=always
RestartSec=5
LimitNOFILE=65535

Environment=HTTP_PORT=443
Environment=WS_PORT=8081
Environment=BASE_DOMAIN=xpo.sh
Environment=SCHEME=https
Environment=REGION=eu1
Environment=RUST_LOG=xpo_server=info
Environment=ACME_ENABLED=true
Environment=ACME_STAGING=false

EnvironmentFile=-/etc/xpo/server.env

AmbientCapabilities=CAP_NET_BIND_SERVICE
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/xpo
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
